# robotpi/DEBIAN/postinst
#!/bin/bash

systemctl enable robotpi
systemctl start robotpi
