# robotpi/DEBIAN/prerm
#!/bin/bash

systemctl stop robotpi
systemctl disable robotpi
